# src/main/resources/cassandra/single-version-testing/supervisord.conf
[supervisord]
nodaemon=true
autorestart=false

[program:sshd]
command=/usr/sbin/sshd -D

[program:cassandra]
command=/bin/bash -c "source /usr/bin/set_env && /usr/local/bin/cassandra-clusternode"

[program:cqlsh]
command=/bin/bash -c "source /usr/bin/set_env && /cassandra/cqlsh_init.sh"
stopasgroup=true
stopsignal=KILL
startretries=10

[group:upfuzz_cassandra]
programs=cassandra,cqlsh
